# build/minecraft-server-wrapper/contrib/minecraft-server-wrapper.service
[Unit]
Description=minecraft-server-wrapper
After=network.target

[Service]
User=minecraft
ExecStart=/usr/bin/minecraft-server-wrapper
StandardOutput=inherit
StandardError=inherit
Restart=always
EnvironmentFile=/etc/defaults/minecraft-server-wrapper

[Install]
WantedBy=multi-user.target
